# conferences/RedisGears/supervisord.conf
[unix_http_server]
file = /run/supervisord.sock

[supervisord]
nodaemon = true
pidfile = /run/supervisord.pid

logfile = /dev/stdout
logfile_maxbytes = 0

[supervisorctl]
serverurl = unix:///run/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:ttyd]
user = student
directory = /src/redisu
command = /usr/local/bin/ttyd --port 8081 /bin/bash
autostart = true
autorestart = true

stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:nginx]
# To make Supervisor happy run nginx in foreground ('daemon off;'):
#command = /run/wait-for-it.sh -t 10 localhost:8080 -- /usr/sbin/nginx -g 'daemon off;'
command = /usr/sbin/nginx -g 'daemon off;'
directory = /var/lib/nginx/html
user = nginx
autostart = true
autorestart = true

[program:redis]
command = redis-server --loadmodule /var/opt/redislabs/lib/modules/rejson.so --loadmodule /var/opt/redislabs/lib/modules/redistimeseries.so --loadmodule /var/opt/redislabs/lib/modules/redisearch.so --loadmodule /var/opt/redislabs/lib/modules/redisbloom.so --loadmodule /var/opt/redislabs/lib/modules/redisgears.so --loadmodule /var/opt/redislabs/lib/modules/redisgraph.so PythonHomeDir /opt/redislabs/lib/modules/python3
user = redis 
autostart = true
autorestart = true
directory = /redis5data
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:redis6]
command = redis6-server --port 6370
user = redis 
autostart = true
autorestart = true
directory = /redis6data
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0